# assets/config/supervisor/6-cdp-proxy.conf
[program:cdp-proxy]
command=python3 /opt/cdp_proxy.py
directory=/opt
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
