# assets/config/supervisor/7-api-server.conf
[program:api-server]
command=python3 /opt/api_server.py
priority=100
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
